# mysqlBackup.cron
0 0 * * * sh /home/sunmnet/rule/cloud-v4.1.0/mysqlbak/mysql_dump_script.sh
